# examples/system/etc/udev/rules.d/80-i915.rules
ACTION=="add", SUBSYSTEM=="module", KERNEL=="i915", RUN+="/usr/bin/proptest -M i915 CONNECTOR_NUMBER connector ID_FROM_SECOND_COLUMN 1"
